# systemd/user/reminders.service
[Unit]
Description=Personal Reminders

[Service]
Type=oneshot
ExecStart=/usr/bin/notify-send "💧 ପିୟ ଜଳ! ୧୦ ପାଦ ଚାଲା"
