# systemd/user/syncdots.service
[Unit]
Description=Sync Dots

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/update-dots.sh
